# make_snapshot.service
# Makes a system snapshot

[Unit]
Description=System Snapshotting

[Service]
Type=oneshot
ExecStart=/snapshots/snapshot.sh
